# profile.d/pantheon-filechooser-module.sh
export GTK_MODULES="$GTK_MODULES:pantheon-filechooser-module"
