# profile.d/xlib-skip-argb-visuals.sh
# emacs workaround
export XLIB_SKIP_ARGB_VISUALS=1
